# deploy/systemd/lumina-web.service
###############################################################################
# lumina-web.service
# Blossom Lumina — 단일 통합 WEB 서비스 (Asset Mgmt + Chat + Dashboard)
# 하나의 systemd 유닛이 두 개의 gunicorn 프로세스를 동반 기동/정지한다.
#   - 127.0.0.1:8000 : Lumina Dashboard (python3.6)
#   - 127.0.0.1:8001 : 자산관리 + 채팅 (python3.11)
# 설치 경로: /usr/lib/systemd/system/lumina-web.service
###############################################################################

[Unit]
Description=Lumina Web (Asset Mgmt + Chat + Dashboard)
Documentation=https://blossom.local/docs/lumina-web
After=network-online.target mariadb.service
Wants=network-online.target nginx.service

[Service]
Type=simple
User=root
Group=root

# ── 환경 변수 ────────────────────────────────────────────
EnvironmentFile=-/etc/blossom/lumina/secure.env
Environment=FLASK_ENV=production
Environment=PYTHONDONTWRITEBYTECODE=1

# ── 두 gunicorn을 한 wrapper에서 동반 기동 ───────────────
ExecStart=/usr/local/bin/lumina-web-start.sh
ExecReload=/bin/kill -HUP $MAINPID

# ── 재시작 / 종료 ───────────────────────────────────────
KillMode=mixed
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

# ── 로그 ─────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lumina-web

[Install]
WantedBy=multi-user.target
