# deploy/systemd/lumina-ap.service
###############################################################################
# lumina-ap.service
# Blossom Lumina — AP 서버 (Agent 데이터 수신/처리/DB 적재)
# 설치 경로: /usr/lib/systemd/system/lumina-ap.service
###############################################################################

[Unit]
Description=Blossom Lumina AP Server (Agent Data Receiver/Processor)
Documentation=https://blossom.local/docs/lumina-ap
After=network-online.target mariadb.service
Wants=network-online.target
# DB가 같은 서버에 있는 경우에만 활성화
# Requires=mariadb.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
# ── 실행 계정 (root 실행 금지) ────────────────────────────
Type=notify
User=lumina-ap
Group=lumina-ap

# ── 환경 변수 (비밀값 분리 로드) ─────────────────────────
EnvironmentFile=-/etc/blossom/lumina/secure.env

# ── 실행 명령 ────────────────────────────────────────────
ExecStart=/usr/bin/python3 /opt/blossom/lumina/ap/wsgi.py \
    --config /etc/blossom/lumina/ap.conf
ExecReload=/bin/kill -HUP $MAINPID

# ── 작업 디렉터리 ────────────────────────────────────────
WorkingDirectory=/opt/blossom/lumina/ap

# ── 재시작 정책 ──────────────────────────────────────────
Restart=on-failure
RestartSec=10

# ── 타임아웃 ─────────────────────────────────────────────
TimeoutStartSec=30
TimeoutStopSec=30
WatchdogSec=120

# ── PID 파일 ─────────────────────────────────────────────
PIDFile=/run/blossom/lumina/ap.pid
RuntimeDirectory=blossom/lumina
RuntimeDirectoryMode=0755

# ── 로그 ─────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lumina-ap

# ━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━
# 보안 강화 (systemd 샌드박스)
# ━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━

# 파일시스템 보호
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/blossom/lumina/ap
ReadWritePaths=/var/log/blossom/lumina/ap
ReadWritePaths=/run/blossom/lumina
ReadOnlyPaths=/etc/blossom/lumina
ReadOnlyPaths=/opt/blossom/lumina

# 임시 파일 격리
PrivateTmp=yes

# 커널 보호
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# 네트워크 (수신 :5100 + DB :3306 연결만 필요)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# 권한 제한
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=
# 5100 포트는 1024 이상이므로 CAP_NET_BIND_SERVICE 불필요하나,
# 향후 포트 변경 대비 명시적으로 허용

# 장치 접근 금지
PrivateDevices=yes

# 시스템 콜 제한
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@mount @reboot @swap @clock @debug @obsolete @raw-io

# 메모리 보호
MemoryDenyWriteExecute=yes

# 사용자 네임스페이스
PrivateUsers=yes

# ulimit
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
