# deploy/systemd/dropins/nginx.service.d/lumina.conf
###############################################################################
# nginx.service drop-in (Blossom Lumina)
# 설치 경로: /etc/systemd/system/nginx.service.d/lumina.conf
# 효과: lumina-web.service stop/restart 시 nginx도 함께 stop/restart
###############################################################################
[Unit]
PartOf=lumina-web.service
